# tvh-streamer-fixed.service
[Unit]
Description=TVHeadend Streamer GPU Transcoder
Documentation=https://github.com/adis992/tvheadend-streamer-panel
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=noname
Group=noname
WorkingDirectory=/home/noname/tvheadend-streamer-panel
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /home/noname/tvheadend-streamer-panel/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=process
TimeoutSec=30

# Security settings (relaxed for home directory access)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/noname/tvheadend-streamer-panel/streams
ReadWritePaths=/home/noname/tvheadend-streamer-panel/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tvh-streamer

[Install]
WantedBy=multi-user.target
